# tools/gunicorn.conf
[program:gunicorn]
command=gunicorn wsgi -c /home/proj/directory/api/gunicorn.conf.py
directory=/home/proj/directory/api/
user=nobody
autostart=true
autorestart=true
redirect_stderr=true